# roles/druss/templates/xbindkeys.service
[Unit]
Description=XBindKeys

[Service]
ExecStart=/usr/bin/xbindkeys

[Install]
WantedBy=default.target